# pkg/moonbow-backend/etc/supervisor.d/moonbow.ini
[program:moonbow-server]
command=/usr/bin/python "/opt/moonbow-backend/server.py"
autostart=true
autorestart=true

