# pkg/moonbow-backend/etc/supervisor.d/moonbow.ini
[program:moonbow]
command=/usr/bin/python "/opt/moonbow-backend/moonbowServer.py"
autostart=true
autorestart=true

