# pkg/moonbow-backend/etc/supervisor.d/moonbow.ini
[program:moonbow-server]
command=/usr/bin/python "/opt/moonbow-backend/server.py"
autostart=true
autorestart=true

[program:moonbow-button]
command=/usr/bin/python "/opt/moonbow-backend/buttonListener.py"
autostart=true
autorestart=true
