# pkg/moonbow-backend/etc/supervidor.d/moonbow.ini
[program:moonbow]
command=python /opt/moonbow/moonbowServer.py
autostart=true
autorestart=true

